# docker/supervisord.dev.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10

[program:laravel]
command=/bin/sh -c "cd /var/www/html && composer install --no-interaction && php artisan serve --host=0.0.0.0 --port=8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=20
startsecs=10

[program:vite]
command=/bin/sh -c "cd /var/www/html && npm install && npm run dev -- --host 0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30
startsecs=15

[program:queue]
command=/bin/sh -c "sleep 30 && cd /var/www/html && php artisan queue:work --sleep=3 --tries=3"
stdout_logfile=/var/www/html/storage/logs/queue.log
stderr_logfile=/var/www/html/storage/logs/queue.log
autorestart=true
priority=40
startsecs=30
